# AGENT_FOX_V1B1B/PACKAGE_MANIFEST.json
{
  "all_inline_validation_scripts_compiled": true,
  "authoritative_minimal_hashes": {
    "MINIMAL_LOAD_CAPTURE_RECEIPT.json": "19a76bef2adf4f88ba495ed8abcf48e356b52f628de97eccabaf6a334fbd7c82",
    "MINIMAL_LOAD_PROCESS_SUMMARY.json": "83b37ebcdb6b5882573f74e96f81b42296e4d6a6fc7a2f1e9a9158c8546ae76f",
    "MINIMAL_LOAD_QUALIFICATION.json": "c300ef1e519b0dc0be949d528034ad8548392e0791e544ce70e4a650f7a09026",
    "MINIMAL_LOAD_RESOURCE_SNAPSHOT.json": "5ca895a436b91047b095a97efaaa7de2e23ea35555c9b17664a844826db1c639"
  },
  "authoritative_minimal_load_mission": "ENG-20260722-062243-69C555",
  "automatic_process_changes": false,
  "change_count": 3,
  "collector_self_test_passed": true,
  "collector_sha256": "1dee5e6e0870a75a7a38181095849880f43acd78caf6782ea17cf932e5ddb146",
  "collector_size_bytes": 44200,
  "contract_sha256": "a873d8f81ecf1ad7169bcaee6410823736981f993ce7a0e7740bd0f6a368c226",
  "contract_size_bytes": 2159,
  "desktop_source_modified": false,
  "drop_folder": "AGENT_FOX_V1B1B",
  "exact_generated_output_count": 4,
  "existing_live_source_files_modified": 0,
  "image_generation_activity_detection": "not evaluated non-invasively; user must keep ComfyUI idle",
  "k_accessed": false,
  "minimal_load_evidence_modified": false,
  "mission_id": "ENG-20260722-064558-2C21EB",
  "network_connections_initiated": false,
  "normal_loaded_precondition_required": true,
  "package_files": [
    "PLAN.json",
    "pc_loaded_resource_capture_v1.py",
    "PC_LOADED_RESOURCE_CAPTURE_CONTRACT_V1.json",
    "README_V1B1B.md",
    "README_FIRST.txt",
    "PACKAGE_MANIFEST.json",
    "PACKAGE_BUILD_RECEIPT.json"
  ],
  "processes_changed": false,
  "readme_sha256": "a240b80f1b69836b8fc959e50b1345e69da468af10219494c64c69cef2b1e158",
  "readme_size_bytes": 1013,
  "registry_writes": false,
  "required_components": [
    "foxai_desktop",
    "model_runtime_or_loopback_8080",
    "comfyui_loopback_8188",
    "webui_loopback_8765"
  ],
  "schema": "foxai.agent_fox.technical_core.v1b1b.package_manifest.v1",
  "services_changed": false,
  "shared_adapter_modified": false,
  "startup_items_changed": false,
  "title": "Agent Fox Technical Core V1B-1B — Normal Loaded Resource Baseline",
  "validation_count": 8,
  "webui_source_modified": false
}
